# scripts/udev/100-udev-usbserv-launch.rules
ACTION=="add", KERNEL=="ttyGS0", RUN+="/home/crtx/usb_serv/scripts/launch.sh"
