# scripts/udev/100-udev-usbserv-stop.rules
ACTION=="remove", KERNEL=="ttyGS0", RUN+="/home/crtx/usb_serv/scripts/stop.sh"
